# .codex/copilot_bridge/ops/systemd/copilot-bridge.service
[Unit]
Description=Copilot Bridge API (Express wrapper for Copilot CLI)
After=network.target

[Service]
Type=simple
User=runner
Group=runner
WorkingDirectory=/srv/copilot-bridge
EnvironmentFile=/etc/copilot-bridge/env
# Honor PATH from EnvironmentFile so we can scope Node 22 if needed:
ExecStart=/usr/bin/env node .codex/copilot_bridge/bridge/server.js
Restart=on-failure
RestartSec=3s
KillMode=process
TimeoutStopSec=10

# Optional hardening (tune per distro)
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=full
ProtectHome=read-only

[Install]
WantedBy=multi-user.target
